# supervisord.conf
[program:uwsgi]
command=uwsgi --ini /home/workspace/api/uwsgi.ini
directory=/home/workspace/api
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true